# Daemon/daemonWallpaper.service
[Unit]
Description=Service to change wallpaper

[Service]
Type=simple

ExecStart=/bin/bash /opt/wallpaperChanger.sh